# deploy/systemd/qiita-compute-orchestrator.service.d/umask.conf
# UMask=0007 ensures qiita-orch creates per-step attempt dirs (the
# nested input/output/logs under each work-ticket workspace) mode 2770
# instead of systemd's default 0022 → 0755. qiita-job (the SLURM
# job-execution identity) then writes outputs via group-write access
# through the qiita-pipeline group. See "Known gap — nested workspace
# dir perms" in docs/runbooks/first-deploy.md for the failure mode this
# fixes.
[Service]
UMask=0007
