# deploy/systemd/qiita-data-plane@.service
[Unit]
Description=Qiita Data Plane (instance %i)
After=network.target postgresql.service

[Service]
Type=simple
User=qiita-data
ExecStart=/opt/qiita/data-plane/qiita-data-plane
Restart=on-failure
Environment=LISTEN_ADDR=127.0.0.1:%i
EnvironmentFile=/etc/qiita/data-plane.env

[Install]
WantedBy=multi-user.target
