# deploy/systemd/qiita-data-plane@.service
[Unit]
Description=Qiita Data Plane (instance %i)
After=network.target postgresql.service

[Service]
Type=simple
User=qiita-data
# DuckDB writes its extension cache under $HOME; qiita-data's passwd home is
# /dev/null (non-login svc account), so set a real one. Per-instance dir
# (%i) avoids concurrent-write races between simultaneous instances when
# horizontally scaled (qiita-data-plane@50051, @50052, …).
# Nested StateDirectory (subpath after the first segment) needs systemd ≥235;
# Rocky 10 ships systemd 256 so this is fine. On older distros, flatten the
# path to e.g. StateDirectory=qiita-data-%i and update HOME= accordingly.
StateDirectory=qiita-data/%i
Environment=HOME=/var/lib/qiita-data/%i
EnvironmentFile=/etc/qiita/data-plane.env
# AFTER EnvironmentFile= on purpose. systemd applies Environment= and
# EnvironmentFile= in the order they appear, last writer wins — so with this line
# ABOVE the env file, a LISTEN_ADDR set in the SHARED /etc/qiita/data-plane.env
# would override the per-instance value for EVERY instance, and they would all
# fight over one port. The instance specifier is the port and must always win.
Environment=LISTEN_ADDR=127.0.0.1:%i
ExecStart=/opt/qiita/data-plane/qiita-data-plane
Restart=on-failure

[Install]
WantedBy=multi-user.target
