# deploy/systemd/qiita-control-plane.service.d/umask.conf
# UMask=0007 ensures the CP runner (qiita-api) creates per-work-ticket
# workspace subdirs mode 2770 instead of systemd's default 0022 → 0755.
# Combined with qiita-orch's membership in qiita-pipeline (see
# docs/runbooks/first-deploy.md §0.1), this lets the orchestrator
# mkdir input/output/logs under those subdirs via group-write access.
# See "Known gap — nested workspace dir perms" in first-deploy.md for
# the failure mode this fixes.
[Service]
UMask=0007
